# etc/supervisor/conf.d/zabbix-agent.conf
[program:zabbix-agent]
command=/usr/sbin/zabbix_agentd -f
user=zabbix
directory=/etc/zabbix
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
